# etc/supervisord.d/sshd.ini
[program:sshd]
command=/opt/qnib/sshd/bin/start.sh
redirect_stderr=true
stdout_syslog=true
