# etc/supervisord.d/yarn-nodemanager.ini
[program:yarn-nodemanager]
priority=940
command=/opt/qnib/yarn/nodemanager/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
