# etc/supervisord.d/yarn-resourcemanager.ini
[program:yarn-resourcemanager]
priority=930
command=/opt/qnib/yarn/resourcemanager/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
